# may-2026/devops/quickstart/deploy/systemd/alchemyst-api.service
[Unit]
Description=Alchemyst JSON API gateway
Requires=iii-engine.service
After=network-online.target iii-engine.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=alchemyst
Group=alchemyst
WorkingDirectory=/opt/alchemyst-hiring/may-2026/devops/quickstart/app/api-gateway
EnvironmentFile=/etc/alchemyst/stack.env
EnvironmentFile=/etc/alchemyst/api-gateway.env
ExecStartPre=/usr/bin/timeout 30s bash -c "until nc -z 127.0.0.1 49134; do sleep 1; done"
ExecStart=/usr/local/bin/run-api-gateway.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
