# may-2026/devops/quickstart/deploy/systemd/alchemyst-inference-worker.service
[Unit]
Description=Alchemyst inference worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alchemyst
Group=alchemyst
WorkingDirectory=/opt/alchemyst-hiring/may-2026/devops/quickstart/workers/inference-worker
EnvironmentFile=/etc/alchemyst/stack.env
EnvironmentFile=/etc/alchemyst/inference-worker.env
ExecStart=/usr/local/bin/run-inference-worker.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
